# misc/felixweb.conf
description     "Felix Web Server"
author          "Erick Tryzelaar <erick.tryzelaar@gmail.com>"

start on (net-device-up
          and local-filesystems
          and runlevel [2345])
stop on runlevel [016]

#expect daemon
#respawn

env HOME=/home/felixweb
env LD_LIBRARY_PATH=/usr/local/lib/felix/felix-latest/host/lib/rtl 
umask 007

script
        privbind -u felixweb /usr/local/lib/felix/felix-latest/host/bin/flx_web --port=80 --root=/usr/local/lib/felix/felix-latest
end script

